# packaging/prerm
#!/bin/sh

deb-systemd-invoke stop grafsy.service
deb-systemd-helper disable grafsy.service
